# csphere/csphere/files/units/csphere-backup.timer
[Unit]
Description=Run Csphere Backup Every 60 minutes

[Timer]
OnBootSec=30
OnCalendar=*-*-* *:00:00
